# root/etc/s6-overlay/s6-rc.d/svc-opendkim/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec /usr/sbin/opendkim -p /run/opendkim/opendkim.sock -P /run/opendkim/pid -d $MYDOMAIN -k /config/dkim/default.private -u abc
